# paulsha_cortex/deploy/templates/manager.timer.tmpl
[Unit]
Description=Deprecated __INSTANCE__ persona manager timer (paulsha-cortex)

[Timer]
OnBootSec=2min
OnUnitActiveSec=300
Unit=__INSTANCE__-manager.service

[Install]
WantedBy=timers.target
